# app-deploy.service
[Unit]
Description=Start cuidabem-api docker
After=network.target docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/home/app/uninter/cuidabem
ExecStart=/usr/bin/docker compose -f docker-compose.yml -f docker-compose.prod.yml up -d --build
RemainAfterExit=true
User=app
Group=app

[Install]
WantedBy=multi-user.target
